# infra/systemd/fairway-field.service
# /etc/systemd/system/fairway-field.service
#
# Polls ESPN for each upcoming tournament whose field hasn't been
# published yet, and stamps `tournaments.field_published_at` the
# first time ESPN returns a non-empty competitors collection.
#
# What gets updated:
#   - `golfers` (new field members get inserted via espn_id upsert)
#   - `scores`  (zero-score rows seeded so the picks UI can filter
#                the dropdown by JOIN scores ON tournament_id)
#   - `tournaments.field_published_at` (stamped on first non-empty
#                response — unlocks the picks UI for that tournament)
#
# Cheap when no field is pending — runFieldSync() short-circuits after
# one SELECT count if there are no candidate tournaments inside the
# 14-day horizon.

[Unit]
Description=Fairway Fantasy — poll ESPN for upcoming tournament fields
After=fairway-fantasy.service network-online.target
Wants=network-online.target
Requires=fairway-fantasy.service

[Service]
Type=oneshot
User=greg
EnvironmentFile=/opt/fairway-fantasy/.env.local

ExecStart=/usr/bin/curl -fsS --max-time 60 \
  -X POST \
  -H "Authorization: Bearer ${CRON_SECRET}" \
  http://127.0.0.1:3000/api/sync-field

SuccessExitStatus=0
